# root/etc/s6-overlay/s6-rc.d/init-dolphin-config/up
/etc/s6-overlay/s6-rc.d/init-dolphin-config/init.sh